# package/DEBIAN/postinst
#!/bin/bash
insserv /etc/init.d/kaos 
systemctl enable kaos
/etc/init.d/kaos start

