# kaos/package/DEBIAN/prerm
#!/bin/bash
/etc/init.d/kaos stop
systemctl disable kaos
update-rc.d kaos disable
